# config/fifo.conf
[program:*fifo*]
command=sudo mkfifo test.fifo
directory=/tmp
autostart=true
autorestart=false
stderr_logfile=/var/log/supervisor/fifo.log
stdout_logfile=/var/log/supervisor/fifo.log
